# flask/crontab
* * * * * root /bin/bash /app/update_date.sh
# Mandatory blank line