# flask/crontab
* 12,14 * * * root curl -v http://127.0.0.1/update
# Mandatory blank line