# flask/crontab
* * * * * /app/update_date.sh
# Mandatory blank line